# prodigy/crontab
* * * * * . /etc/profile; (cd /app && python -m prodigy db-out dataset > /app/data/export/dataset.jsonl) > /var/log/cron.log 2>&1
